# deploy/starbot-webgui.service
# Systemd service file for Starbot WebGUI
# Place at: /etc/systemd/system/starbot-webgui.service
#
# Setup (run once):
#   sudo useradd -r -s /bin/false starbot
#   sudo chown -R starbot:starbot /var/www/sites/stella/starbot.cloud
#
# Enable: sudo systemctl enable starbot-webgui
# Start:  sudo systemctl start starbot-webgui
# Status: sudo systemctl status starbot-webgui
# Logs:   sudo journalctl -u starbot-webgui -f

[Unit]
Description=Starbot WebGUI Service (Next.js)
After=network.target starbot-api.service

[Service]
Type=simple
User=starbot
Group=starbot
WorkingDirectory=/var/www/sites/stella/starbot.cloud/Starbot_WebGUI
Environment="NODE_ENV=production"
Environment="PORT=3001"
Environment="NEXT_PUBLIC_API_URL=https://starbot.cloud/v1"

# Start command (Next.js production server)
ExecStart=/usr/bin/node /var/www/sites/stella/starbot.cloud/Starbot_WebGUI/.next/standalone/server.js

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=starbot-webgui

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
PrivateDevices=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
